# infra/spark/deploy.master/etc/profile.d/spark.sh
export SPARK_CONF_DIR=/root/spark
export PATH=$PATH:${SPARK_CONF_DIR}/bin
