# infra/spark/deploy.master/etc/profile.d/spark.sh
export SPARK_HOME=/root/spark
export SPARK_CONF_DIR=${SPARK_HOME}/conf
export PATH=$PATH:${SPARK_HOME}/bin
