# infra/spark/deploy.generic/etc/profile.d/hadoop.sh
export PATH=$PATH:/root/ephemeral-hdfs/bin
