# infra/ec2/deploy.generic/etc/profile.d/prompt.sh
export PS1="\e[36m({{cluster_name}})\e[0m \u@\h:\e[32m\W\e[0m \$ "